# pcap.conf
[inet_http_server]
port=127.0.0.1:23456

[supervisorctl]
serverurl=http://127.0.0.1:23456

[supervisord]
logfile = /dev/null
loglevel = error
identifier = pcap

[program-default]
envFiles = /pcap.env
numprocs = 1
stdout_logfile = /dev/stdout
stderr_logfile = /dev/stderr
directory = /
stopsignal = TERM
stopwaitsecs = 10
exitcodes = 0
startsecs = 0

[program:gcsfuse]
envFiles=/env/gcsfuse.env
command = /scripts/start_gcsfuse
process_name = gcsfuse
priorirty = 0

[program:gcsdir]
envFiles=/env/gcsdir.env
command = /scripts/create_pcap_dir
process_name = gcsdir
autorestart = false
depends_on = gcsfuse
priorirty = 1

[program:pcapfsn]
envFiles=/env/pcapfsn.env
command = /scripts/start_pcapfsn 
process_name = pcapfsn
depends_on = gcsdir
priorirty = 2

[program:tcpdumpw]
envFiles=/env/rt.env,/env/tcpdumpw.env
command = /scripts/start_tcpdumpw
process_name = tcpdumpw
depends_on = pcapfsn
priorirty = 3

[group:pcap_sidecar]
programs=gcsfuse,gcsdir,pcapfsn,tcpdumpw
